# deploy/supervisor/yoza-queue.conf
; ══════════════════════════════════════════════════════════════════════
; Supervisor config — Queue Worker (Yoza Rent Car)
;
; Cara pakai di server (Ubuntu/Debian + Supervisor terinstall):
;   1. sudo cp deploy/supervisor/yoza-queue.conf /etc/supervisor/conf.d/
;   2. Sesuaikan path "/var/www/yoza-rent-car" & user "www-data" di bawah
;      dengan path/user yang sebenarnya di server Anda.
;   3. sudo supervisorctl reread
;   4. sudo supervisorctl update
;   5. sudo supervisorctl start yoza-queue:*
;
; Setiap kali deploy ulang / ada kode baru:
;   php artisan queue:restart
;   (supervisor otomatis menghidupkan ulang worker dengan kode terbaru)
; ══════════════════════════════════════════════════════════════════════

[program:yoza-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/yoza-rent-car/artisan queue:work redis --sleep=3 --tries=3 --backoff=10 --max-time=3600
directory=/var/www/yoza-rent-car
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/yoza-rent-car/storage/logs/queue-worker.log
stopwaitsecs=3600
