# deploy/supervisor/yoza-reverb.conf
; ══════════════════════════════════════════════════════════════════════
; Supervisor config — Reverb WebSocket Server (Yoza Rent Car)
;
; Reverb adalah proses panjang (long-running) untuk fitur chat realtime.
; Wajib dijalankan terus-menerus di belakang web server — bukan lewat
; "php artisan serve".
;
; Cara pakai di server:
;   1. sudo cp deploy/supervisor/yoza-reverb.conf /etc/supervisor/conf.d/
;   2. Sesuaikan path "/var/www/yoza-rent-car" & user "www-data" di bawah.
;   3. Pastikan port REVERB_PORT (default 8080) di .env sudah dibuka di
;      firewall, dan reverse-proxy (nginx) sudah diarahkan ke port itu
;      untuk path WebSocket (lihat dokumentasi Reverb: proxy WSS → ws).
;   4. sudo supervisorctl reread
;   5. sudo supervisorctl update
;   6. sudo supervisorctl start yoza-reverb:*
; ══════════════════════════════════════════════════════════════════════

[program:yoza-reverb]
process_name=%(program_name)s
command=php /var/www/yoza-rent-car/artisan reverb:start
directory=/var/www/yoza-rent-car
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/yoza-rent-car/storage/logs/reverb.log
stopwaitsecs=30
